# roles/vagrant/templates/vboxclient.service.j2
[Unit]
Description=VirtualBox X Window System guest services
Requires=x11.target

[Service]
ExecStart=/usr/bin/VBoxClient --{{ item }} --nodaemon

[Install]
WantedBy=default.target
